# provisioning/resources/init/snowplow_mini.service
[Unit]
Description=Snowplow Mini Services
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/home/ubuntu/snowplow
ExecStartPre=/usr/sbin/service docker restart
ExecStartPre=docker-compose down
ExecStart=docker-compose up -d
ExecStop=docker-compose down

[Install]
WantedBy=multi-user.target