# provisioning/roles/sp_mini_7_setup_init/files/snowplow-report.timer
[Unit]
Description=Starts Snowplow Report timer sending mini status metrics
Requires=snowlow-report.service

[Timer]
Unit=snowlow-report.service
OnCalendar=*:0/01:00

[Install]
WantedBy=timers.target